# gunicorn/django-gunicorn.socket
[Unit]
Description=Django gunicorn socket
[Socket]
ListenStream=/run/django-gunicorn.sock
[Install]
WantedBy=sockets.target
